# provisioning/templates/etc/supervisord.d/todo.ini.j2
[program:todo]
command=/home/vagrant/app/app.py 
autostart=true
autorestart=true
user=vagrant
redirect_stderr=true
stdout_logfile=/var/log/todo.log